# devops/app/supervisor/conf.d/collectstatic.conf
[program:collectstatic]
directory=/app
command=python3 manage.py collectstatic --noinput
redirect_stderr=true
priority=1